# docker-events-exporter.service
[Unit]
Description=Docker Events Prometheus Exporter
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
# Adjust path to where you install the binary
ExecStart=/usr/local/bin/docker-events-exporter
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
